# launchd/com.tradingview.scanner.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
  Scanner A — Premarket Gappers
  Runs Mon-Fri at 15:00 Berlin (= 09:00 ET, 30 min before US market open)

  Install:
    1. Replace __HOME__ with your home directory (e.g. /Users/yourname)
    2. Replace __PYTHON_PATH__ with: dirname $(which python3)
    3. cp this file ~/Library/LaunchAgents/
    4. launchctl load ~/Library/LaunchAgents/com.tradingview.scanner.plist
-->
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.tradingview.scanner</string>
    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>__HOME__/TradingViewScanner/daily_scanner.sh</string>
    </array>
    <key>StartCalendarInterval</key>
    <array>
        <dict><key>Weekday</key><integer>1</integer><key>Hour</key><integer>15</integer><key>Minute</key><integer>0</integer></dict>
        <dict><key>Weekday</key><integer>2</integer><key>Hour</key><integer>15</integer><key>Minute</key><integer>0</integer></dict>
        <dict><key>Weekday</key><integer>3</integer><key>Hour</key><integer>15</integer><key>Minute</key><integer>0</integer></dict>
        <dict><key>Weekday</key><integer>4</integer><key>Hour</key><integer>15</integer><key>Minute</key><integer>0</integer></dict>
        <dict><key>Weekday</key><integer>5</integer><key>Hour</key><integer>15</integer><key>Minute</key><integer>0</integer></dict>
    </array>
    <key>StandardOutPath</key>
    <string>__HOME__/TradingViewScanner/launchd_stdout.log</string>
    <key>StandardErrorPath</key>
    <string>__HOME__/TradingViewScanner/launchd_stderr.log</string>
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin:__PYTHON_PATH__</string>
    </dict>
</dict>
</plist>
